# test/automated/crawl/supervisord.conf
[supervisord]
nodaemon=true

[program:virtual-web]
directory=/usr/mywi
command=node test/automated/crawl/virtual-web/index.js
autorestart=false
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
startsecs=5

[program:tests]
directory=/usr/mywi
command=npm run test-crawler
autorestart=false
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
startsecs=0